# wifi_mouse_linux/files/wifi_mouse.service
[Unit]
Description=Vulnerable remote mouse software.

[Service]
ExecStart=/usr/bin/mouseserver

[Install]
WantedBy=multi-user.target